# lib/itamae/plugin/recipe/bedrock_server/templates/etc/systemd/system/bedrock_server.socket.erb
[Unit]
BindsTo=bedrock_server.service

[Socket]
ListenFIFO=%t/bedrock_server
